# supervisord.conf
[supervisord]
nodaemon=true

[program: intellisense]
command = uwsgi /usr/src/app/uwsgi.ini
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
user=root
